# 15min/A8/dashboard-v3-codex/systemd/polymarket-live.service
[Unit]
Description=Polymarket Live Trader (10 USDC capped)
After=network-online.target

[Service]
Type=simple
WorkingDirectory=/root/projects/polymarket-arb/src
ExecStart=/usr/bin/python3 -u live_trader.py
StandardOutput=append:/root/projects/polymarket-arb/output/live_trader.log
StandardError=append:/root/projects/polymarket-arb/output/live_trader.log
Restart=always
RestartSec=10

[Install]
WantedBy=default.target
